# adguard/rootfs/etc/s6-overlay/s6-rc.d/adguard/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# AdGuard Home Plus
# Runs the AdGuard Home server
# ==============================================================================
declare -a options

bashio::log.info "Starting AdGuard Home server..."

options+=(--port 45158)
options+=(--host 127.0.0.1)
options+=(--work-dir /data/adguard)
options+=(--config /config/AdGuardHome.yaml)
options+=(--no-check-update)

if bashio::debug; then
  option+=(--verbose)
fi

cd /opt/AdGuardHome || exit
exec ./AdGuardHome "${options[@]}"
