# adguard/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: AdGuard Home
# Runs the NGINX daemon
# ==============================================================================

# Wait for adguard to become available
bashio::net.wait_for 45158 127.0.0.1 900

bashio::log.info "Starting NGINX..."
exec nginx
